# config/logwatch-ai.cron
#!/bin/bash
# Cron job configuration for logwatch-ai
# Place this file at: /etc/cron.d/logwatch-ai
# Then run: sudo systemctl restart cron

# Minimal environment for cron execution
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# EDIT Logwatch-AI execution directory
EXEDIR=/opt/logwatch-with-ai

# Execute daily at 2:00 AM
0 2 * * * root cd "$EXEDIR" && "$EXEDIR"/.venv/bin/python src/main.py >> /var/log/logwatch-ai-cron.log 2>&1
